# root/conreqroot/root/etc/s6-overlay/s6-rc.d/init-conreq/up
#!/command/execlineb
# shellcheck disable=SC1000-SC9999

/etc/s6-overlay/s6-rc.d/init-conreq/script