# root/etc/services.d/conreq/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/conreq || exit 1

exec \
    s6-setuidgid abc daphne conreq.asgi:application \
    -b 0.0.0.0
